# xmenu/.local/service/xclickroot-xmenu/run
#!/bin/sh
cd
exec xclickroot -r xmenu-open.sh xmenu-root.sh
