# i3/.local/service/i3/run
#!/bin/sh
cd
export EXIT=$(ps -o ppid= -p $(ps -o ppid= -p $$))
exec i3
