# stow/x11vnc/.local/service/x11vnc/run
#!/bin/sh
exec x11vnc -noxdamage -o /dev/null -passwd thv7847r
