# launchar/.local/service/xclickroot/run
#!/bin/sh
cd
exec xclickroot -m launchar
