# stow/acpid/.local/service/acpid/run
#!/bin/sh
[ -r conf ] && . ./conf
exec acpid -fS ${OPTS:=-c $HOME/.config/acpid/events}
